# crontab
PATH=/srv/tw/bin_:/usr/local/bin:/usr/bin:/bin:/usr/local/games:/usr/games

# clean up worse ranks every night
24 3    * * *   racecleanranks
