# systemd/user/twrecordd.service
[Service]
Type=simple
ExecStart=/srv/tw/race/recordd.py
Restart=always

[Install]
WantedBy=default.target
